# trunk/x11-misc/ya-layout/files/etc/udev/rules.d/88-autofs.rules
SUBSYSTEM=="block", RUN+="/bin/bash /usr/share/ya-layout/autofs-all"

