# trunk/net-misc/pspacer/files/usr/ppatch/dev-libs/libnl/preinst/pspd
#!/sbin/runscript

depend() {
	need net
}

start() {
	ebegin "Starting pspd"
	/usr/sbin/pspd -d
	eend $?
}

stop() {
	ebegin "Stopping pspd"
	killall /usr/sbin/pspd
	eend $?
}
